# systemd/mailsync.service
[Unit]
Description=Check Mail
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=/home/vasudev/.local/bin/mail-sync
